# docker/etc/profile.d/nxgit.sh
#!/bin/bash
export NXGIT_CUSTOM=/data/nxgit
