# docker/etc/s6/nxgit/run
#!/bin/bash
[[ -f ./setup ]] && source ./setup

pushd /app/nxgit > /dev/null
    exec su-exec $USER /app/nxgit/nxgit web
popd
